# config/.config/autostart/pulseeffects.desktop
[Desktop Entry]
Version=1.0
Name=PulseEffects
Icon=pulseeffects
Exec=pulseeffects
Terminal=false
Type=Application
Hidden=true
